# packaging/hyprpilot.service
[Unit]
Description=hyprpilot — Tauri overlay daemon for agent-driven workflows
Documentation=https://github.com/hyprpilot/hyprpilot
# Bind to the graphical session so the unit only runs while the
# captain has a desktop / Wayland session active. PartOf =
# graphical-session.target stops the service automatically when the
# session exits.
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# --hidden boots the overlay invisible; first user-visible map is the
# Hyprland keybind / tray click / `hyprpilot ctl overlay show`.
ExecStart=/usr/bin/hyprpilot daemon --hidden
Restart=on-failure
RestartSec=5
# Hardening — the daemon is captain-scoped (unix socket under
# $XDG_RUNTIME_DIR), so no system-wide access is needed.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t %h/.config/hyprpilot %h/.local/state/hyprpilot %h/.local/share/hyprpilot

[Install]
WantedBy=graphical-session.target
